# sigma-rules/windows/registry/win-persistence-run-key.yml
title: Registry Run Key Persistence Mechanism
id: d5e6f7a8-b9c0-1234-5678-901234ef0123
status: stable
level: medium
description: Detects modifications to registry Run/RunOnce keys for establishing persistence.
author: Shakil Md. Rezwanul Bari
date: 2024/01/15
references:
    - https://attack.mitre.org/techniques/T1547/001/
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        TargetObject|contains:
            - '\\CurrentVersion\\Run\\'
            - '\\CurrentVersion\\RunOnce\\'
            - '\\CurrentVersion\\RunServices\\'
    filter_legitimate:
        Details|contains:
            - 'Microsoft'
            - 'Google'
            - 'Adobe'
            - 'Windows Defender'
    condition: selection and not filter_legitimate
falsepositives:
    - Legitimate software installation
tags:
    - attack.persistence
    - attack.t1547.001
